# applications/aox/applications/positioning/debconf/prerm
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop uic-positioning.service
    systemctl disable uic-positioning.service
fi