# applications/openthread_border_router/debconf/services/otbr-agent.service.d/override.conf
[Unit]
PartOf=uic-otbr.service
BindsTo=uic-cpcd.service
After=uic-cpcd.service
