# applications/dev_ui/dev_gui/debconf/scripts/prerm
#!/bin/sh
systemctl stop uic-dev-gui-api.service
systemctl stop uic-dev-gui-server.service
systemctl disable uic-dev-gui-api.service
systemctl disable uic-dev-gui-server.service
