# applications/aox/applications/aoxpc/debconf/prerm
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop uic-aoxpc.service
    systemctl disable uic-aoxpc.service
fi