# applications/cpc-hci-bridge/debconf/prerm
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop uic-cpc-hci-bridge.service
    systemctl disable uic-cpc-hci-bridge.service
fi