# applications/cpcd/debconf/services/uic-cpcd.service.in
[Unit]
Description=UIC CPC Daemon
ConditionPathExists=/var/lib/uic/cpcd-binding-key
Wants=uic-cpcd-bind.service
After=uic-cpcd-bind.service

[Service]
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_BINDIR}/cpcd -c ${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_SYSCONFDIR}/cpcd.conf --key /var/lib/uic/cpcd-binding-key
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
KillMode=process
User=uic
Group=uic
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
