# applications/openthread_border_router/debconf/services/uic-otbr.service.in
[Unit]
Description=UIC Openthread Border Router

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/uic-otbr
ExecStart=${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_DATAROOTDIR}/otbr/script/server

[Install]
WantedBy=multi-user.target

