# applications/zigpc/debconf/prerm.in
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop uic-zigpc.service
    systemctl disable uic-zigpc.service
fi
