# applications/openthread_border_router/debconf/services/otbr-web.service.d/override.conf
[Unit]
PartOf=uic-otbr.service
Requires=otbr-agent.service

[Service]
ExecStartPre=bash -c "test -e /run/openthread-wpan0.sock || sleep 5"
