# applications/cpcd/debconf/prerm.in
#!/bin/sh
if type "systemctl" > /dev/null; then
    systemctl stop uic-cpcd.service
    systemctl disable uic-cpcd.service
fi
