# applications/dev_ui/dev_gui/debconf/services/uic-dev-gui-client.service
[Unit]
Description=UIC User interface
After=network.target mosquitto.service

[Service]
ExecStart=/usr/bin/dev-gui-client
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
KillMode=process
User=uic
Group=uic

Environment=HTTPS=true
Environment=SSL_CRT_FILE=/var/lib/dev-gui/cert.pem
Environment=SSL_KEY_FILE=/var/lib/dev-gui/key.pem

[Install]
WantedBy=multi-user.target
