# applications/zigbeed/debconf/services/uic-zigbeed.service.in
[Unit]
Description=Zigbeed service
StartLimitIntervalSec=0
BindsTo=uic-cpcd.service uic-zigbeed-socat.service
After=uic-cpcd.service uic-zigbeed-socat.service

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/stdbuf -o0 ${CMAKE_INSTALL_PREFIX}/${CMAKE_INSTALL_BINDIR}/zigbeed
ExecStop=/bin/kill -WINCH ${MAINPID}
PIDFile=/run/zigbeed.pid
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=zigbeed

[Install]
WantedBy=multi-user.target
