# applications/openthread_border_router/debconf/prerm.in
#!/bin/sh
if type "systemctl" >/dev/null; then
    systemctl stop otbr-agent.service
    systemctl disable otbr-agent.service
fi
